# cluster/playbooks/roles/music/templates/sendspin.service
[Unit]
Description=Run sendspin in daemon mode

[Service]
ExecStart=/home/ongy/.local/bin/sendspin daemon
